# supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:bridge-pointofsale]
command=/opt/een/app/point_of_sale/bridge-pointofsale
directory=/opt/een/app/point_of_sale
autostart=true
autorestart=true
startretries=3
startsecs=10
stopsignal=TERM
stderr_logfile=/opt/een/data/var/log/point_of_sale/bridge-pointofsale-err.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
stdout_logfile=/opt/een/data/var/log/point_of_sale/bridge-pointofsale-out.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=10