# supervisord.conf
[supervisord]
nodaemon=true
logfile=/opt/een/var/log/point_of_sale/supervisord.log
pidfile=/tmp/supervisord.pid
user=root

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:bridge-devices-pos]
command=/opt/een/app/point_of_sale/bridge-devices-pos
directory=/opt/een/app/point_of_sale
user=root
autostart=true
autorestart=true
startretries=3
redirect_stderr=true
stdout_logfile=/opt/een/var/log/point_of_sale/bridge-devices-pos.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
environment=PATH="/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"